# backup_current/systemd/daily-health-check.timer
[Unit]
Description=Timer for Daily AI System Health Check

[Timer]
OnCalendar=daily
Persistent=true
Unit=daily-health-check.service

[Install]
WantedBy=timers.target
